# config/systemd/user/beets-web.service
[unit]
Description=Beets Web

[Service]
ExecStart=/usr/bin/beet web

[Install]
WantedBy=default.target

